# espcloud-sync/rootfs/etc/services.d/espcloud-sync/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the espcloud-sync service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

bashio::log.info "Starting Espcloud-sync..."

log_level=$(bashio::config 'log_level')

if [ -z "$log_level" ]
then
    export LOGLEVEL="${log_level}"
fi


api_token=$(bashio::config 'api_token')
host=$(bashio::config 'host')

export API_TOKEN="${api_token}"
export HOST="${host}"


exec python3 /usr/bin/sync.py
